# etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
	ExecStart=-/sbin/agetty --autologin "orangepi" %I
